# ansible/roles/openbsd-email/templates/clamsmtpd.rc
#!/bin/ksh
daemon="/usr/local/sbin/clamsmtpd"
daemon_flags="-f /etc/clamsmtpd.conf"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
